# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:postgres]
command=gosu postgres postgres -D "%(ENV_PGDATA)s"
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/supervisor/postgres.log
stderr_logfile=/var/log/supervisor/postgres.err.log
stopsignal=INT

[program:redis]
command=redis-server --dir "%(ENV_REDIS_DATA_DIR)s" --bind 0.0.0.0 --protected-mode no
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/supervisor/redis.log
stderr_logfile=/var/log/supervisor/redis.err.log

[program:kafka]
command=/entrypoint-kafka.sh
autostart=true
autorestart=true
priority=20
stdout_logfile=/var/log/supervisor/kafka.log
stderr_logfile=/var/log/supervisor/kafka.err.log
startsecs=10

[program:app]
command=/entrypoint-app.sh
autostart=true
autorestart=true
priority=30
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
